# resources/etc/services.d/nginx/run
#!/bin/sh
/usr/sbin/nginx
echo "nginx running"