# docker/resources/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:php7-fpm]
command=/usr/local/sbin/php-fpm
